# worker/stackhealth-worker.service
[Unit]
Description=Stack Health Video Merge Worker
After=network.target
Wants=network.target

[Service]
Type=simple
User=stackhealth
Group=stackhealth
WorkingDirectory=/opt/stackhealth-worker
EnvironmentFile=/opt/stackhealth-worker/.env
ExecStart=/opt/stackhealth-worker/venv/bin/python worker.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stackhealth-worker

# 안전한 ffmpeg 실행을 위한 권한
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
